# files/etc/crontabs/root
0 1 * * 0 echo "" > /tmp/ssrplus.log
0 2 * * * /usr/share/shadowsocksr/subscribe.sh
0 5 * * 0 /usr/share/shadowsocksr/update.sh
*/15  * * * * /etc/lz/CheckLoger
