# self_files/etc/crontabs/root
*/15  * * * * /etc/lz/CheckLoger
16 5 * * * sleep 5 && touch /etc/banner && reboot
0 2 * * * /usr/share/shadowsocksr/ssrplusupdate.sh >>/tmp/ssrplus.log 2>&1
